# moodle/crontab
0 1 * * * php /var/www/html/admin/cli/cron.php
